# debian/softgred.prerm
if [ -x "/etc/init.d/softgred" ]; then
	invoke-rc.d softgred stop || exit $?
fi

