# debian/softgred.postrm
if [ "$1" = "purge" ] ; then
	update-rc.d softgred remove >/dev/null
fi

